# devel/etc/catkin/profile.d/06-ctr-nuke.sh
#!/usr/bin/env sh
unset CATKIN_TEST_RESULTS_DIR
unset ROS_TEST_RESULTS_DIR
